# backup/112_queue_db_backup.service
[Unit]
Description=15-112 Queue DB Backup

[Service]
Type=simple
ExecStart=/home/queue-dev/15-112-Queue/backup/manual_backup
